# server/cron/cron.daily/roadsituation.cron
#! /bin/sh
cd /srv/www/htdocs/roadsituation/server/scripts/
php /srv/www/htdocs/roadsituation/server/scripts/set_gati.php >> /var/log/roadsituation/gati.log
php /srv/www/htdocs/roadsituation/server/scripts/set_gibdd.php >> /var/log/roadsituation/gibdd.log
php /srv/www/htdocs/roadsituation/server/scripts/set_speedcamonline.php >> /var/log/roadsituation/speedcamonline.log
php /srv/www/htdocs/roadsituation/server/scripts/set_transportmos.php >> /var/log/roadsituation/transportmos.log
php /srv/www/htdocs/roadsituation/server/scripts/set_i.centr.by.php >> /var/log/roadsituation/set_i.centr.by.log
php /srv/www/htdocs/roadsituation/server/scripts/set_saratov.gov.ru.php >> /var/log/roadsituation/set_saratov.gov.ru.log
php /srv/www/htdocs/roadsituation/server/scripts/set_ufacity.info.php >> /var/log/roadsituation/set_ufacity.info.log
php /srv/www/htdocs/roadsituation/server/scripts/set_ekaterinburg.php >> /var/log/roadsituation/set_ekaterinburg.log
