# server/cron/cron.hourly/roadsituation.cron
#! /bin/sh
curator --config /git/roadsituation/elasticsearch/curator/config.yaml /git/roadsituation/elasticsearch/curator/roadsituation.yaml >> /var/log/curator/roadsituation.log